# host/overlay/usr/lib/systemd/system/omnilab-nvidia-cdi.service
[Unit]
Description=Generate NVIDIA Container Device Interface for podman
Documentation=https://github.com/dhworg/omnilab
ConditionPathExists=!/etc/cdi/nvidia.yaml
ConditionPathExists=/usr/bin/nvidia-ctk
After=local-fs.target
Before=podman.service podman.socket multi-user.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/mkdir -p /etc/cdi
ExecStart=/usr/bin/nvidia-ctk cdi generate --output=/etc/cdi/nvidia.yaml
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
